# secure-channel/deploy/systemd/dl-ids.service
[Unit]
Description=Darklock IDS (Identity Discovery Service)
After=network.target
Wants=network-online.target
Documentation=https://github.com/darklock/secure-channel

[Service]
Type=simple
User=darklock
Group=darklock
WorkingDirectory=/opt/darklock/secure-channel/services/ids
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=5
StartLimitBurst=5
StartLimitIntervalSec=60

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ReadWritePaths=/opt/darklock/secure-channel/services/ids/data
ReadOnlyPaths=/opt/darklock/secure-channel/services/ids

# Environment
EnvironmentFile=-/etc/darklock/ids.env
Environment=NODE_ENV=production
Environment=PORT=4100

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dl-ids

[Install]
WantedBy=multi-user.target
